# usr/lib/crontab
0,15,30,45 * * * * /usr/lib/atrun
0,10,20,30,40,50 * * * * /etc/dmesg - >> /usr/adm/messages
